# supervisord.conf
[supervisord]
nodaemon=true

[program:hub]
command=/opt/bin/entry_point.sh
user=seluser

[program:node]
command=/opt/bin/node_entry_point.sh
user=seluser
